# files/etc/update-motd.d/40-about-puppet
#!/bin/bash

echo "
This machine is under PUPPET control.
To update:
  cd /etc/puppet
  sudo git pull
  sudo puppet apply
But seriously - don't do that without cloning and tesing eh?
"
